# scripts/rh-adapter/rh-adapter.service
[Unit]
Description=RunningHub Adapter for NewAPI Billing
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/rh-adapter
ExecStart=/usr/bin/node server.mjs
Restart=always
RestartSec=5
EnvironmentFile=/opt/rh-adapter/.env
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
